# node-terminal-server.service
[Unit]
Description=Node Terminal Server

[Service]
Type=notify
ExecStart=/usr/bin/node INSTALLDIR/server.js

